# profile.d/userprofile.sh
#!/bin/bash

# Aliases

alias egrep='egrep --color=auto'
alias fgrep='fgrep --color=auto'
alias grep='grep --color=auto'
alias ls='ls --color=auto'
alias l='ls -CF'
alias la='ls -A'
alias ll='ls -alhF'
alias ..='cd ..'
alias ports='netstat -tulanp'
alias guid="uuidgen | tr '[:upper:]' '[:lower:]'"
